# docker/rootfs/etc/services.d/go2rtc/run
#!/command/with-contenv bashio
# ==============================================================================
# Runs go2rtc
# ==============================================================================
# shellcheck shell=bash

bashio::log.info 'Starting go2rtc...'

declare config_path
declare binary_path

config_path="/config"

binary_path="/usr/local/bin/go2rtc"
if [[ -x "${config_path}/go2rtc" ]]; then
  binary_path="${config_path}/go2rtc"
  bashio::log.warning "Using go2rtc binary from '${binary_path}'"
fi

# set cwd for go2rtc (for config file, Hass integration, etc)
cd "${config_path}" || bashio::exit.nok "Could not change working directory"

exec "${binary_path}"
